# training_container/training.cron
00 18 * * 5  cd /app/python && /usr/local/bin/python /app/python/training.py >> /var/log/cron.log 2>&1
#* * * * * cd /app/python && /usr/local/bin/python /app/python/training.py >> /var/log/cron.log 2>&1
